# config/ai.openclaw.hook-briefing.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.openclaw.hook-briefing</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>HOOK_REPO_PATH/scripts/morning-briefing.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>HOOK_USER_HOME</string>
        <key>HOOK_DIR</key>
        <string>HOOK_REPO_PATH</string>
        <key>HOOK_SLACK_CHANNEL</key>
        <string>#hook</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>HOOK_REPO_PATH</string>

    <key>StandardOutPath</key>
    <string>HOOK_USER_HOME/.openclaw/logs/hook/briefing-launchd.log</string>

    <key>StandardErrorPath</key>
    <string>HOOK_USER_HOME/.openclaw/logs/hook/briefing-launchd.err</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
